# profile.d/mail.sh
export MAIL=/var/vmail/
export MAIL=/var/vmail/~
export MAIL=/var/vmail/
export MAIL=~/Maildir:LAYOUT=fs:DIRNAME=/var/vmail/
